# files/etc/uci-defaults/99-cosmos-default
# Beware! This script will be in /rom/etc/uci-defaults/ as part of the image.



# Uncomment lines to apply:

# if left commented, LAN will be set to dhcp client
# lan_ip_address="192.168.1.1"

# enable cosmos on boot
cosmos=1

if [ -n "$(ifconfig -a | grep 'wlan0')" ]; then
  wlan_name="Cosmos-Cloud-$(ifconfig -a | grep -m 1 'wlan0' | tail -c 8 | head -c 5 | tr -d : )"
  wlan_password=$(tr -dc A-Za-z0-9 </dev/urandom | head -c 13; echo)
fi
root_password=$(tr -dc A-Za-z0-9 </dev/urandom | head -c 13; echo)

#################

# log potential errors
exec >/tmp/99-cosmos-default.log 2>&1

# disable checking for upgrades (it fails anyway due to rootfs size limitations)
uci set attendedsysupgrade.client.login_check_for_upgrades='0'

# disable Password-Login for SSH = only allow authorized_keys
uci set dropbear.main.PasswordAuth='off'
uci set dropbear.main.RootPasswordAuth='off'
uci commit dropbear

if [ -f "/boot/.env" ]; then source /boot/.env; fi
if [ "${SSH_KEY}" ]; then echo "${SSH_KEY}" >> /etc/dropbear/authorized_keys ; fi

echo "PASSWD=$root_password" >> /opt/cosmos-config/.env

# Configure LAN
# More options: https://openwrt.org/docs/guide-user/base-system/basic-networking
if [ -n "$lan_ip_address" ]; then
  echo "set static IP $lan_ip_address"
  uci set network.lan.ipaddr='$lan_ip_address'
  uci commit network  
else
  echo "set IP to dhcp-client"
  uci set network.lan.proto='dhcp'
  uci set network.lan11=interface
  uci set network.lan11.name='staticFallback'
  uci set network.lan11.device='@lan'
  uci set network.lan11.proto='static'
  uci set network.lan11.ipaddr='192.168.0.11/24'
  uci commit network
  uci add_list firewall.@zone[0].network='lan11'
  uci commit firewall
fi

# Configure WLAN
# More options: https://openwrt.org/docs/guide-user/network/wifi/basic#wi-fi_interfaces
if [ -n "$wlan_name" -a -n "$wlan_password" -a ${#wlan_password} -ge 8 ]; then
  echo "set Wifi-Settings"
  uci set wireless.@wifi-device[0].disabled='0'
  uci set wireless.@wifi-device[0].htmode='HT40'
  uci set wireless.@wifi-device[0].band='2g'
  uci set wireless.@wifi-device[0].channel='11'
  uci set wireless.@wifi-device[0].legacy_rates='1'
  uci set wireless.@wifi-device[0].country='CH'
  uci set wireless.@wifi-device[0].cell_density='0'
  uci set wireless.@wifi-iface[0].disabled='0'
  uci set wireless.@wifi-iface[0].encryption='psk2'
  uci set wireless.@wifi-iface[0].multicast_to_unicast_all='1'
  uci set wireless.@wifi-iface[0].ieee80211w='1'
  uci set wireless.@wifi-iface[0].ocv='0'
  uci set wireless.@wifi-iface[0].wpa_disable_eapol_key_retries='1'
  uci set wireless.@wifi-iface[0].disassoc_low_ack='0'

  uci set wireless.@wifi-iface[0].ssid="$wlan_name"
  uci set wireless.@wifi-iface[0].key="$wlan_password"
  uci commit wireless
fi

# setup cosmos
if [ -n "$cosmos" ]; then
  
  #set hostname & Timezone
  uci set system.@system[0].hostname="cosmos-cloud-$(ifconfig -a | grep -m 1 'eth0' | tail -c 8 | head -c 5 | tr -d : )"
  uci set system.@system[0].zonename='Europe/Zurich'
  uci set system.@system[0].timezone='CET-1CEST,M3.5.0,M10.5.0/3'
  uci commit system

  #store hostname in env-var
  echo "HOSTNAME=0.0.0.0" >> /opt/cosmos-config/.env

  #switch luci to free ports
  echo "setup luci to free ports"
  uci del uhttpd.main.listen_http
  uci add_list uhttpd.main.listen_http='0.0.0.0:8080'
  uci del uhttpd.main.listen_https
  uci commit uhttpd

  #add network rules for cosmos-Networks
  echo "setup firewall rules for cosmos"
  uci add firewall zone
  uci set firewall.@zone[-1].name='cosmos'
  uci set firewall.@zone[-1].input='ACCEPT'
  uci set firewall.@zone[-1].output='ACCEPT'
  uci set firewall.@zone[-1].forward='ACCEPT'
  uci add_list firewall.@zone[-1].subnet='172.16.0.0/12'
  uci add firewall forwarding
  uci set firewall.@forwarding[-1].src='cosmos'
  uci set firewall.@forwarding[-1].dest='lan'
  uci add firewall forwarding
  uci set firewall.@forwarding[-1].src='cosmos'
  uci set firewall.@forwarding[-1].dest='wan'
  uci add firewall forwarding
  uci set firewall.@forwarding[-1].src='lan'
  uci set firewall.@forwarding[-1].dest='cosmos'
<<<<<<< HEAD
  #uci add firewall zone
  #uci set firewall.@zone[-1].name='constellation'
  #uci set firewall.@zone[-1].input='ACCEPT'
  #uci set firewall.@zone[-1].output='ACCEPT'
  #uci set firewall.@zone[-1].forward='REJECT'
  #uci add_list firewall.@zone[-1].subnet='192.168.201.0/24'
=======
  uci add firewall zone
  uci set firewall.@zone[-1].name='constellation'
  uci set firewall.@zone[-1].input='ACCEPT'
  uci set firewall.@zone[-1].output='ACCEPT'
  uci set firewall.@zone[-1].forward='ACCEPT'
  uci add_list firewall.@zone[-1].subnet='192.168.201.0/24'
  uci add firewall forwarding
  uci set firewall.@forwarding[-1].src='constellation'
  uci set firewall.@forwarding[-1].dest='lan'
  uci add firewall forwarding
  uci set firewall.@forwarding[-1].src='constellation'
  uci set firewall.@forwarding[-1].dest='wan'

  uci set network.lan.metric='10'
  #uci set network.lan11.metric='20'
  uci set network.wan.metric='30'

  uci set mwan3.lan=interface
  uci set mwan3.lan.initial_state='online'
  uci set mwan3.lan.family='ipv4'
  uci add_list mwan3.lan.track_ip='1.0.0.1'
  uci add_list mwan3.lan.track_ip='1.1.1.1'
  uci add_list mwan3.lan.track_ip='208.67.220.222'
  uci add_list mwan3.lan.track_ip='208.67.222.222'
  uci set mwan3.lan.track_method='ping'
  uci set mwan3.lan.reliability='2'
  uci set mwan3.lan.count='1'
  uci set mwan3.lan.size='56'
  uci set mwan3.lan.max_ttl='60'
  uci set mwan3.lan.timeout='4'
  uci set mwan3.lan.interval='10'
  uci set mwan3.lan.failure_interval='5'
  uci set mwan3.lan.recovery_interval='5'
  uci set mwan3.lan.down='5'
  uci set mwan3.lan.up='5'
  uci set mwan3.lan.enabled='1'

  #uci set mwan3.lan11=interface
  #uci set mwan3.lan11.initial_state='online'
  #uci set mwan3.lan11.family='ipv4'
  #uci add_list mwan3.lan11.track_ip='1.0.0.1'
  #uci add_list mwan3.lan11.track_ip='1.1.1.1'
  #uci add_list mwan3.lan11.track_ip='208.67.220.222'
  #uci add_list mwan3.lan11.track_ip='208.67.222.222'
  #uci set mwan3.lan11.track_method='ping'
  #uci set mwan3.lan11.reliability='2'
  #uci set mwan3.lan11.count='1'
  #uci set mwan3.lan11.size='56'
  #uci set mwan3.lan11.max_ttl='60'
  #uci set mwan3.lan11.timeout='4'
  #uci set mwan3.lan11.interval='10'
  #uci set mwan3.lan11.failure_interval='5'
  #uci set mwan3.lan11.recovery_interval='5'
  #uci set mwan3.lan11.down='5'
  #uci set mwan3.lan11.up='5'
  #uci set mwan3.lan11.enabled='1'

  uci del mwan3.wanb
  uci del mwan3.wanb6
  uci del mwan3.wanb_m1_w2
  uci del mwan3.wanb_m1_w3
  uci del mwan3.wanb_m2_w2
  uci del mwan3.wanb6_m1_w2
  uci del mwan3.wanb6_m1_w3
  uci del mwan3.wanb6_m2_w2
  uci del mwan3.wan_m1_w3
  uci del mwan3.wan_m2_w3
  uci del mwan3.wan6_m1_w3
  uci del mwan3.wan6_m2_w3
  uci set mwan3.lan_m1_w1=member
  uci set mwan3.lan_m1_w1.interface='lan'
  uci set mwan3.lan_m1_w1.metric='1'
  uci set mwan3.lan_m1_w1.weight='1'
  #uci set mwan3.lan11_m2_w1=member
  #uci set mwan3.lan11_m2_w1.interface='lan11'
  #uci set mwan3.lan11_m2_w1.metric='2'
  #uci set mwan3.lan11_m2_w1.weight='1'
  uci set mwan3.wan_m3_w1=member
  uci set mwan3.wan_m3_w1.interface='wan'
  uci set mwan3.wan_m3_w1.metric='3'
  uci set mwan3.wan_m3_w1.weight='1'
  uci del mwan3.wanb_only
  uci del mwan3.wan_wanb
  uci del mwan3.wanb_wan
  uci del mwan3.balanced.use_member
  uci add_list mwan3.balanced.use_member='wan_m1_w3'
  uci add_list mwan3.balanced.use_member='wan6_m1_w3'
  uci set mwan3.default=policy
  uci add_list mwan3.default.use_member='lan_m1_w1'
  #uci add_list mwan3.default.use_member='lan11_m2_w1'
  uci add_list mwan3.default.use_member='wan_m3_w1'
  uci set mwan3.default.last_resort='default'
  uci del mwan3.https
  uci set mwan3.default_rule_v4.use_policy='default'
  uci set mwan3.default_rule_v6.use_policy='default'

  if [ -e /dev/cdc-wdm0 ]; then
    uci set network.wwan=interface
    uci set network.wwan.proto='qmi'
    #uci set network.globals.packet_steering='1'
    uci set network.wwan.device='/dev/cdc-wdm0'
    uci set network.wwan.auth='none'
    uci set network.wwan.pdptype='ipv4'
    uci set network.wwan.metric='40'

    uci add_list firewall.@zone[1].network='wwan'

    uci set mwan3.wwan=interface
    uci set mwan3.wwan.initial_state='online'
    uci set mwan3.wwan.family='ipv4'
    uci add_list mwan3.wwan.track_ip='1.0.0.1'
    uci add_list mwan3.wwan.track_ip='1.1.1.1'
    uci add_list mwan3.wwan.track_ip='208.67.220.222'
    uci add_list mwan3.wwan.track_ip='208.67.222.222'
    uci set mwan3.wwan.track_method='ping'
    uci set mwan3.wwan.reliability='2'
    uci set mwan3.wwan.count='1'
    uci set mwan3.wwan.size='56'
    uci set mwan3.wwan.max_ttl='60'
    uci set mwan3.wwan.timeout='4'
    uci set mwan3.wwan.interval='10'
    uci set mwan3.wwan.failure_interval='5'
    uci set mwan3.wwan.recovery_interval='5'
    uci set mwan3.wwan.down='5'
    uci set mwan3.wwan.up='5'
    uci set mwan3.wwan_m4_w1=member
    uci set mwan3.wwan_m4_w1.interface='wwan'
    uci set mwan3.wwan_m4_w1.metric='4'
    uci set mwan3.wwan_m4_w1.weight='1'
    uci add_list mwan3.balanced.use_member='wwan_m4_w1'
    uci add_list mwan3.default.use_member='wwan_m4_w1'
    uci set mwan3.wwan.enabled='1'
  fi

  uci commit network
  uci commit mwan3
>>>>>>> 3ed946a (add constellation firewall rules)
  uci commit firewall
fi

echo "All done!"
exit 0